# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-extra-ui.target
[Unit]
Description=RudyRainOS Stage 3: Extra UI Elements
After=rudy-profiles.target
