# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-extensions.target
[Unit]
Description=RudyRainOS Stage 4: Extensions and Background Services
After=rudy-extra-ui.target
