# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-ui.target
[Unit]
Description=RudyRainOS Stage 1: Core UI
After=graphical-session.target
