# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-profiles.target
[Unit]
Description=RudyRainOS Stage 2: User Profiles
After=rudy-ui.target
