# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-apps.target
[Unit]
Description=RudyRainOS Stage 5: Applications
After=rudy-extensions.target
