# Meta-Packages/rudyraindesktop/files/etc/systemd/system/rudy-everything-else.target
[Unit]
Description=RudyRainOS Stage 6: Remaining System Services
After=rudy-apps.target
