# Meta-Packages/rudyraindesktop/debian/postinst
#!/bin/bash
set -e

# Debian postinst for rudyraindesktop
# - safe guards for build/chroot vs real installed system
# - only run user session actions when a real user and DBus session exist
# - create /var/log/execas.log and sample /etc/execasers if missing

REAL_USER="${SUDO_USER:-$USER}"
# Try to find a real user home even if SUDO_USER empty (e.g., during install)
if [ -n "$REAL_USER" ]; then
    USER_HOME="$(getent passwd "$REAL_USER" | cut -d: -f6 || true)"
fi

# detect if we are likely inside a container/chroot/build environment:
init_name="$(basename "$(readlink -f /proc/1/exe 2>/dev/null || true)")"
SKIP_SYSTEM=false
if [ "$init_name" != "systemd" ]; then
    SKIP_SYSTEM=true
fi

# Helper to safely run a command as the real user if possible
run_as_user() {
    local uid="$1"; shift
    if [ -z "$REAL_USER" ] || [ -z "$uid" ] || [ ! -d "$USER_HOME" ]; then
        return 0
    fi
    # prefer sudo if available
    if command -v sudo >/dev/null 2>&1; then
        sudo -u "$REAL_USER" XDG_RUNTIME_DIR="/run/user/$uid" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$uid/bus" "$@"
    else
        # try su fallback
        su - "$REAL_USER" -c "XDG_RUNTIME_DIR=/run/user/$uid DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/$uid/bus $*"
    fi
}

# Create /var/log/execas.log if packaging expects it
if [ ! -e /var/log/execas.log ]; then
    touch /var/log/execas.log || true
    chown root:root /var/log/execas.log 2>/dev/null || true
    chmod 0644 /var/log/execas.log 2>/dev/null || true
fi

# Ensure /etc/execasers exists (packaging should install sample); create conservative sample if absent
if [ ! -e /etc/execasers ]; then
    cat > /etc/execasers <<'EOF' || true
# /etc/execasers - sample execas rules
# Format:
#  !caller allowed_target1 allowed_target2 ...
# Example:
!alice all
!sysadmin root
EOF
    chown root:root /etc/execasers 2>/dev/null || true
    chmod 0644 /etc/execasers 2>/dev/null || true
fi

# User-specific configuration: only attempt when we have a real user and a runtime bus
if [ -n "$REAL_USER" ] && [ -n "$USER_HOME" ]; then
    REAL_UID="$(id -u "$REAL_USER" 2>/dev/null || true)"
    if [ -n "$REAL_UID" ] && [ -d "/run/user/$REAL_UID" ] && [ -S "/run/user/$REAL_UID/bus" ]; then
        # Set wallpaper if asset exists
        if [ -f /usr/share/rudyraindesktop/wallpapers/default.jpg ]; then
            run_as_user "$REAL_UID" gsettings set org.gnome.desktop.background picture-uri "file:///usr/share/rudyraindesktop/wallpapers/default.jpg" || true
            run_as_user "$REAL_UID" gsettings set org.gnome.desktop.background picture-uri-dark "file:///usr/share/rudyraindesktop/wallpapers/default.jpg" || true
        fi

        # Import terminal profile if present
        if [ -f /usr/share/rudyraindesktop/terminal-profiles/RudyRainOS.json ]; then
            run_as_user "$REAL_UID" dconf load /org/gnome/terminal/legacy/profiles:/ < /usr/share/rudyraindesktop/terminal-profiles/RudyRainOS.json || true
        fi

        # Enable GNOME extensions safely (only if extension binary is available)
        if command -v gnome-extensions >/dev/null 2>&1; then
            run_as_user "$REAL_UID" gnome-extensions enable user-theme@gnome-shell-extensions.gcampax.github.com || true
            run_as_user "$REAL_UID" gnome-extensions enable appindicatorsupport@rgcjonas.gmail.com || true
            run_as_user "$REAL_UID" gnome-extensions enable dash-to-dock@micxgx.gmail.com || true
            run_as_user "$REAL_UID" gnome-extensions enable blur-my-shell@aunetx || true
            run_as_user "$REAL_UID" gnome-extensions enable rounded-window-corners@yilozt || true
        fi

        # User-level gsettings tweaks
        run_as_user "$REAL_UID" gsettings set org.gnome.desktop.interface font-name 'Inter 11' || true
        run_as_user "$REAL_UID" gsettings set org.gnome.desktop.interface monospace-font-name 'JetBrains Mono 11' || true
        run_as_user "$REAL_UID" gsettings set org.gnome.desktop.peripherals.touchpad tap-to-click true || true
        run_as_user "$REAL_UID" gsettings set org.gnome.desktop.peripherals.touchpad natural-scroll true || true
    fi
else
    # No real user found; skip user-specific customizations
    :
fi

# System-level actions: only when we appear to be on a real installed system (systemd PID 1)
if [ "$SKIP_SYSTEM" = false ]; then
    # Update grub if present
    if command -v update-grub >/dev/null 2>&1; then
        update-grub || true
    fi

    # Install plymouth alternatives if files exist
    if [ -f /usr/share/plymouth/themes/rudyrain/rudyrain.plymouth ]; then
        update-alternatives --install /usr/share/plymouth/themes/default.plymouth default.plymouth /usr/share/plymouth/themes/rudyrain/rudyrain.plymouth 100 || true
        update-alternatives --set default.plymouth /usr/share/plymouth/themes/rudyrain/rudyrain.plymouth || true
        update-initramfs -u || true
    fi

    # Disable apport safely if present
    if [ -f /etc/default/apport ]; then
        sed -i 's/enabled=1/enabled=0/' /etc/default/apport || true
        if command -v systemctl >/dev/null 2>&1; then
            systemctl disable apport.service || true
        fi
    fi
fi

# Compile glib schemas if present
if command -v glib-compile-schemas >/dev/null 2>&1 && [ -d /usr/share/glib-2.0/schemas ]; then
    glib-compile-schemas /usr/share/glib-2.0/schemas/ || true
fi

# Ensure grub theme symlink only if file present
if [ -f /usr/share/gnome-shell/theme/RudyRainOS-gdm.css ]; then
    ln -sf /usr/share/gnome-shell/theme/RudyRainOS-gdm.css /usr/share/gnome-shell/theme/gdm3.css || true
    chmod 644 /usr/share/gnome-shell/theme/RudyRainOS-gdm.css || true
fi

# Run any bundled boot overrides generator if present
if [ -x /usr/lib/rudyraindesktop/generate-boot-overrides.sh ]; then
    /usr/lib/rudyraindesktop/generate-boot-overrides.sh || true
fi

# End
exit 0
